# etc/launchd/com.gthnk.librarian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.gthnk.librarian</string>
	<key>Program</key>
	<string>/Users/idm/.virtualenvs/greenthink/bin/journal_rotate.py</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>3</integer>
	</dict>
	<key>EnvironmentVariables</key>
    <dict>
		<key>SETTINGS</key>
		<string>/Users/idm/.gt/greenthink.conf</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/lib/greenthink-library/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/var/lib/greenthink-library/launchd.err.log</string>
</dict>
</plist>
